# onyx-monero-daemon.service
[Unit]
Description=Onyx Monero Mining Daemon
Documentation=https://github.com/oskodiak/onyx-monero-dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/python3 /home/oskodiak/projects/onyx-monero-dashboard/onyx_miner_daemon.py --foreground
WorkingDirectory=/home/oskodiak/projects/onyx-monero-dashboard
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Environment
Environment=PATH=/run/current-system/sw/bin:/usr/bin:/bin
Environment=HOME=/home/oskodiak
Environment=USER=oskodiak

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=onyx-monero-daemon

[Install]
WantedBy=default.target